# backend/directus/deployment/rsac-directus.service.example
[Unit]
Description=RSAC-UP Directus CMS
After=network.target postgresql.service
Wants=network.target

[Service]
Type=simple
User=rsac
Group=rsac
WorkingDirectory=/opt/rsac_website/backend/directus
EnvironmentFile=/opt/rsac_website/backend/directus/.env
ExecStart=/usr/bin/npm start
Restart=on-failure
RestartSec=5
TimeoutStopSec=30
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
